# linux_root/etc/profile.d/qcimage.sh
. /qcimage/shell/git_functions.sh
. /qcimage/shell/raw_imaging.sh
. /qcimage/shell/settings.sh
. /qcimage/shell/main.sh
